# files/etc/init/node-skeleton.conf
description "node-skelton"
author  "richarvey"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

exec sudo -u node sh -c "/usr/bin/node /home/node/app.js >> /mnt/log/node.log 2>&1"
